# packaging/systemd/xdp47-agent.service
[Unit]
Description=XDP47 Edge Agent
After=network-online.target
Wants=network-online.target

[Service]
Environment=HOME=/var/lib/xdp47
EnvironmentFile=-/etc/default/xdp47-agent
ExecStart=/usr/local/bin/xdp47-agent
Restart=always
RestartSec=3
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
